# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#Scan with 5G radio to fix signal icons in luci
iw wlan0 scan &>/dev/null

exit 0

